# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:vncserver]
command=/bin/bash -c "
    rm -f /tmp/.X1-lock /tmp/.X11-unix/X1 2>/dev/null || true;
    mkdir -p /tmp/.X11-unix;
    chmod 1777 /tmp/.X11-unix;
    export HOME=/home/Nouk;
    export USER=Nouk;
    export LOGNAME=Nouk;
    su - Nouk -c '
        export DISPLAY=:1;
        export HOME=/home/Nouk;
        /usr/bin/tigervncserver :1 \
            -geometry 1280x720 \
            -depth 24 \
            -rfbport 5901 \
            -rfbauth /home/Nouk/.vnc/passwd \
            -fg \
            -nolisten tcp \
            -localhost no \
            -SecurityTypes VncAuth \
            -xstartup /home/Nouk/.vnc/xstartup
    '
"
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=5
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc-err.log

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/supervisor/xrdp-sesman.log
stderr_logfile=/var/log/supervisor/xrdp-sesman-err.log

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/var/log/supervisor/xrdp.log
stderr_logfile=/var/log/supervisor/xrdp-err.log

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5901
autostart=true
autorestart=true
priority=25
startsecs=10
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc-err.log
